# .config/systemd/user/pulseeffects.service
[Unit]
Description=Pulseeffects

[Service]
ExecStart=pulseeffects --gapplication-service

[Install]
WantedBy=default.target
